# webodm_cleaner.timer
[Unit]
Description=Running webodm_cleaner.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target